# new-dotfiles/files/home/.config/systemd/user/sway-session-pre.target
[Unit]
Description=sway compositor session
Documentation=man:systemd.special(7)
BindsTo=graphical-session-pre.target

